# deploy/com.agentic-trading-bot.schedule.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd tick for main.py schedule - -run (see docs/DEPLOYMENT_ROADMAP.md).
     Install:  cp deploy/com.agentic-trading-bot.schedule.plist ~/Library/LaunchAgents/
               launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.agentic-trading-bot.schedule.plist
     Remove:   launchctl bootout gui/$(id -u)/com.agentic-trading-bot.schedule
     Known limitation: a LaunchAgent runs only while logged in and awake -
     the 15:30 ET pre-close gate window needs the lid open. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic-trading-bot.schedule</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>/Users/ninkuk/Desktop/agentic-trading-bot/deploy/schedule-tick.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ninkuk/Desktop/agentic-trading-bot</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <!-- stderr must be captured: the leads legs' skip-and-continue warnings
         only exist there (paper-trail-report reads this file). -->
    <key>StandardOutPath</key>
    <string>/Users/ninkuk/Desktop/agentic-trading-bot/schedule.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ninkuk/Desktop/agentic-trading-bot/schedule.log</string>
</dict>
</plist>
